# bat.service
[Unit]
Description=Battery Low Warning

[Service]
Type=simple
ExecStart=/usr/local/bin/bat.sh
