# chwp.service
[Unit]
Description=Change Wallpaper

[Service]
Type=simple
ExecStart=/usr/bin/chwp.sh
